# host/launchd/com.esl.harness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.esl.harness</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/esl-harness</string>
        <string>daemon-start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/esl-harness.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/esl-harness.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/lite-agent-cli</string>
    <!-- Specify UserName if running as a specific user, otherwise it runs as current user -->
    <!-- <key>UserName</key>
    <string>liteagent</string> -->
</dict>
</plist>
